# external_cloud/csv2/single_container/cloud_scheduler/supervisord.conf
[supervisord]
nodaemon=true
 
[program:sshd]
command=/usr/sbin/sshd -D
 
[program:httpd]
command=/bin/bash -c "rm -rf /run/httpd/* /tmp/httpd* && systemctl start httpd -DFOREGROUND"